# resources/provisioning/whmcs/crontab
*/5 * * * * curl http://localhost/whmcs/crons/cron.php