# resources/provisioning/whmcs/crontab
*/5 * * * * php /var/www/html/whmcs/crons/cron.php
